# deploy.cron
* * * * *   root    /root/min.waf/deploy.sh