# deploy.cron
* * * * *   root    /usr/local/min.waf/deploy.sh